# systemd/.config/systemd/user/obsidian-commit.service
[Unit]
Description=Auto-commit Obsidian vault changes

[Service]
Type=oneshot
ExecStart=%h/.local/bin/obsidian-auto-commit.sh
